# violas/etc/systemd/system/ceph.target
[Unit]
Description=All Ceph clusters and services

[Install]
WantedBy=multi-user.target
