# violas/etc/update-motd.d/50-landscape-sysinfo
/usr/share/landscape/landscape-sysinfo.wrapper